# K2/camera_web/camera-hls.service
[Unit]
Description=USB Camera HLS Capture
After=local-fs.target systemd-udev-settle.service
Wants=systemd-udev-settle.service

[Service]
Type=simple
User=z
Group=z
WorkingDirectory=/home/z/camera-web
ExecStart=/usr/bin/bash /home/z/camera-web/camera_hls_runner.sh
Restart=always
RestartSec=2
StandardOutput=append:/home/z/camera-web/logs/ffmpeg-hls.log
StandardError=append:/home/z/camera-web/logs/ffmpeg-hls.log

[Install]
WantedBy=multi-user.target
