# K2/scripts/deepzzz-k2-edge.service
[Unit]
Description=DeepZZZ K2 Edge Web Service
After=network-online.target local-fs.target
Wants=network-online.target
Conflicts=camera-hls.service camera-web.service

[Service]
Type=simple
User=z
Group=z
WorkingDirectory=/home/z
Environment=HOST=0.0.0.0
Environment=PORT=7860
Environment=DATASET_CAPTURE_PORT=7861
Environment=ENABLE_POSE=1
Environment=POSE_MODEL=/home/z/.brdk_models/pose/yolov8n-pose-320.calib20260601.q.onnx
Environment=POSE_PROVIDER=SpaceMITExecutionProvider
Environment=POSE_PROVIDER_FALLBACK=0
Environment=LOW_LIGHT_MODE=off
Environment=LOW_LIGHT_GAMMA=1.35
Environment=LOW_LIGHT_CLAHE_CLIP=2.0
Environment=LOW_LIGHT_DESATURATE=0
Environment=ENABLE_DATASET_CAPTURE=1
Environment=ANALYSIS_WIDTH=320
Environment=ANALYSIS_FPS=12
Environment=VIDEO_DECODER=h264_stcodec
Environment=AUDIO_DEVICE=plughw:CARD=Device,DEV=0
Environment=AUDIO_RATE=16000
Environment=YAMNET_SECONDS=2
Environment=YAMNET_INTERVAL=0.5
Environment=HLS_TIME=1
Environment=HLS_LIST_SIZE=2
ExecStart=/usr/bin/python3 /home/z/main.py
Restart=always
RestartSec=3
StandardOutput=append:/home/z/runtime/logs/app.out.log
StandardError=append:/home/z/runtime/logs/app.err.log

[Install]
WantedBy=multi-user.target
